# deploy/scripts/gc-firewall.service
[Unit]
# Ground Control public-iface firewall — drop tcp:8000 from public NIC (#828).
#
# Defense in depth on top of ADR-026 bearer auth and the `GC_BIND_IP` tailnet-
# only docker-proxy bind. ADR-030 puts the workload on a single Hetzner host
# with a public IPv4 + IPv6, reachable from the internet by default. Even with
# `GC_BIND_IP` constraining docker-proxy to the tailnet IP, a misconfiguration
# that flips it back to 0.0.0.0 would re-expose the API. This unit drops any
# TCP packet to dport 8000 that arrives on the public interface (`enp8s0`),
# regardless of what the docker-proxy is listening on. Tailnet traffic
# arrives via `tailscale0` (or via `lo` when the host talks to itself by its
# tailnet IP) and is unaffected.
#
# Install (one-time, on red-dragon):
#   sudo install -o root -g root -m 644 \
#     deploy/scripts/gc-firewall.service /etc/systemd/system/gc-firewall.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now gc-firewall.service
#
# Verify:
#   sudo iptables -L INPUT -n -v | head    # rule 1 should DROP -i enp8s0 dport 8000
#   sudo ip6tables -L INPUT -n -v | head   # same on IPv6
#
# Disable:
#   sudo systemctl disable --now gc-firewall.service
#   (ExecStop removes the rules cleanly)
Description=Ground Control public-iface firewall — drop tcp:8000 from public NIC (#828, ADR-026 belt-and-braces on top of GC_BIND_IP)
After=network-online.target tailscaled.service docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=true
# Idempotent: -C checks for existing rule, -I prepends only when missing.
ExecStart=/bin/bash -c '/usr/sbin/iptables -C INPUT -i enp8s0 -p tcp --dport 8000 -j DROP 2>/dev/null || /usr/sbin/iptables -I INPUT 1 -i enp8s0 -p tcp --dport 8000 -j DROP'
ExecStart=/bin/bash -c '/usr/sbin/ip6tables -C INPUT -i enp8s0 -p tcp --dport 8000 -j DROP 2>/dev/null || /usr/sbin/ip6tables -I INPUT 1 -i enp8s0 -p tcp --dport 8000 -j DROP'
ExecStop=/bin/bash -c '/usr/sbin/iptables -D INPUT -i enp8s0 -p tcp --dport 8000 -j DROP 2>/dev/null || true'
ExecStop=/bin/bash -c '/usr/sbin/ip6tables -D INPUT -i enp8s0 -p tcp --dport 8000 -j DROP 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
